# dynamic-dns-server.service
[Unit]
Description=dynamic_dns_server Daemon
Wants=network.target
After=syslog.target network.target

[Service]
Type=forking
User=root
Group=root
WorkingDirectory=/home/mush/code
PIDFile=/var/run/dynamic-dns-server.pid
ExecStart=/home/mush/dynamic-dns/dynamic_dns_server.py
StandardOutput=syslog
ExecStop=/bin/kill -15 $MAINPID
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
